# base-file-structure/autodeploy/autodeploy.service
[Unit]
Description=Autodeploy

[Service]
ExecStart=/home/ubuntu/autodeploy/autodeploy /home/ubuntu/autodeploy/autodeploy.yml

[Install]
WantedBy=multi-user.target
